# rb.user/templates/authorized_keys.service.j2
[Unit]
Description=Update authorized_keys each 5 minutes
Wants=authorized_keys.timer

[Service]
ExecStart=/usr/bin/authorized_keys.sh

[Install]
WantedBy=authorized_keys.target
