# rb.pushmetrics-custom-exporter/templates/metrics.timer.j2
[Unit]
Description=Run metrics-{{ pushgateway_custom_job_name }}.sh every twenty seconds

[Timer]
OnCalendar=*:*:00,20

[Install]
WantedBy=timers.target
