# rb.pushmetrics-node-exporter/templates/metrics.service.j2
[Unit]
Description=send metrics every minute
Wants=metrics.timer

[Service]
ExecStart=/usr/bin/metrics.sh

[Install]
WantedBy=metrics.target
